# packaging/systemd/notify-relayd.service
[Unit]
Description=notify-relay host service
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
EnvironmentFile=-%h/.config/notify-relay/notify-relayd.env
ExecStart=%h/.local/bin/notify-relayd --unix %t/notify-relay.sock
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
